# OpenRC/all-ways-egpu-boot-vga-openrc
#!/sbin/openrc-run

command="all-ways-egpu"
command_args="set-boot-vga egpu"
pidfile="/var/run/${RC_SVCNAME}.pid"
name="All-Ways-eGPU Boot VGA Daemon"

description="This daemon is used by the all-ways-egpu script to automatically set the boot-vga flag at startup"

depends() {
    after display-manager
}
